# etc/systemd/system/greetd.service.d/override.conf
[Unit]
After=systemd-user-sessions.service plymouth-quit-wait.service
After=getty@tty1.service
Conflicts=getty@tty1.service
